# root/etc/services.d/uwsgi/run
#!/usr/bin/with-contenv sh

/usr/local/bin/uwsgi /etc/uwsgi/uwsgi.ini
